# scripts/launchd/com.agent-conductor.scout.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gruai.scout</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/yangyang/Repos/gruai/scripts/scheduled-runs.sh</string>
        <string>scout</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/yangyang/Repos/gruai/logs/launchd-scout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/yangyang/Repos/gruai/logs/launchd-scout-error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/yangyang/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
